# smart_notification_router/rootfs/etc/services.d/smart-notification/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when the service fails
# ==============================================================================

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Smart Notification Router service failed, exiting..."
  exec /run/s6/basedir/bin/halt
fi

bashio::log.info "Smart Notification Router service stopped, restarting..."
