# smart_notification_router/rootfs/etc/services.d/smart-notification/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Smart Notification Router..."

# Generate configuration from the add-on options
bashio::log.info "Generating configuration file from add-on options..."

python3 - << "EOL" > /config/notification_config.yaml
import json
import yaml
import os

options_file = '/data/options.json'

if os.path.exists(options_file):
    with open(options_file, 'r') as f:
        options = json.load(f)
        
    config = {
        'audiences': options.get('audiences', {}),
        'severity_levels': options.get('severity_levels', ['low', 'medium', 'high', 'emergency'])
    }
    
    with open('/config/notification_config.yaml', 'w') as f:
        yaml.dump(config, f, default_flow_style=False)
        
    print("Configuration generated successfully")
else:
    print("Options file not found")
EOL

# Start the application
cd /app
exec python3 main.py