# smart_notification_router/rootfs/etc/services.d/smart-notification/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the smart notification service
# ==============================================================================

# Change to the app directory
cd /app || bashio::exit.nok "Could not change to app directory"

# Run the start script
bashio::log.info "Starting Smart Notification Router..."
exec /app/run.sh
