# smart_notification_router/rootfs/etc/services.d/smart-notification/run
#!/bin/bash
# ==============================================================================
# Start the smart-notification service
# ==============================================================================

echo "[INFO] Starting Smart Notification Router service..."

# Check if the run script exists
if [ -f "/app/run.sh" ]; then
  exec /app/run.sh
else
  echo "[ERROR] Run script not found at /app/run.sh"
  exit 1
fi
