# smart_notification_router/rootfs/etc/services.d/smart-notification/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start the smart notification service
# s6-overlay service for Home Assistant Add-on
# ==============================================================================

# Set up logging
echo "[INFO] Starting Smart Notification Router service..."

# Verify script location
if [ ! -f "/app/run.sh" ]; then
    echo "[ERROR] Run script not found at /app/run.sh"
    find /app -type f -name "*.sh" | sort
    find / -name "run.sh" -type f | sort
    exit 1
fi

# Fix line endings if needed
echo "[INFO] Fixing line endings for run.sh..."
sed -i 's/\r$//' /app/run.sh

# Show permissions
echo "[DEBUG] Script permissions:"
ls -la /app/run.sh

# Check if script has a valid shebang line
echo "[INFO] Checking script shebang line:"
head -n 1 /app/run.sh

# Execute with full error handling
echo "[INFO] Executing run.sh"
exec /bin/bash /app/run.sh
