# smart_notification_router/rootfs/etc/services.d/smart-notification/run
#!/usr/bin/with-contenv bash
# Start your service here

echo "[INFO] Starting Smart Notification Router..."
exec python3 /opt/smart_notification/main.py
