# smart_notification_router/rootfs/etc/s6-overlay/s6-rc.d/smart-notification/run
#!/bin/sh

echo "[DEBUG] Starting smart-notification service"
echo "[DEBUG] Current directory: $(pwd)"
echo "[DEBUG] User: $(whoami)"
echo "[DEBUG] PID: $$"
echo "[DEBUG] PPID: $PPID"
echo "[DEBUG] S6 version check:"
/command/s6-version || echo "S6 command not found"

echo "[DEBUG] Checking for application files:"
ls -la /app || echo "App directory not found"

cd /app
echo "[DEBUG] Now in directory: $(pwd)"

echo "[DEBUG] Starting Python application"
python3 main.py