# smart_notification_router/rootfs/etc/s6-overlay/s6-rc.d/smart-notification/run
#!/command/execlineb -P
s6-setuidgid root
cd /app
with-contenv
python3 main.py