# smart_notification/rootfs/etc/services.d/smart-notification/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH=/data/options.json
NOTIFICATION_CONFIG=/config/notification_config.yaml

bashio::log.info "Starting Smart Notification Router..."

# Generate configuration from the add-on options
bashio::log.info "Generating configuration file..."

python3 - << "EOL" > $NOTIFICATION_CONFIG
import json
import yaml
import os

CONFIG_PATH = '/data/options.json'

with open(CONFIG_PATH, 'r') as f:
    config = json.load(f)

notification_config = {
    'audiences': config['audiences'],
    'severity_levels': config['severity_levels']
}

with open('/config/notification_config.yaml', 'w') as f:
    yaml.dump(notification_config, f, default_flow_style=False)
EOL

# Start the application
cd /app
python3 main.py